# talkha_client/rootfs/etc/services.d/talkha-client/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

HA_TOKEN="$(bashio::config 'ha_token')"
CLIENT_ID="$(bashio::config 'client_id')"

mkdir -p /data/.talkha_state /data/.talkhalokal_state

export PATH="/opt/talkha/venv/bin:${PATH}"
export HOMEASSISTANT_TOKEN="${HA_TOKEN}"

cat > /data/.talkha.env <<EOF
# Generated by TalkHa Client add-on
HA_URL=http://supervisor/core
HA_TOKEN=${HA_TOKEN}
EOF

export TALKHA_OPERATOR_URL="http://10.8.10.10:8787"
export TALKHA_CLIENT_ID="${CLIENT_ID}"
export TALKHA_REGISTRATION_TOKEN="talkha-shared-client-token-2026"
export TALKHA_MODE="full"
export TALKHA_ALLOW_MUTATIONS="true"
export TALKHA_LOCAL_DB_MODE="true"
export TALKHA_ADDON_MODE="1"
export TALKHA_LOCAL_DB="1"
export TALKHA_POLL_INTERVAL="10"

exec python3 /opt/talkha/app/agent.py
