# talkha_client/rootfs/etc/services.d/talkha-client/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

HA_URL="$(bashio::config 'ha_url')"
HA_TOKEN="$(bashio::config 'ha_token')"
OPERATOR_URL="$(bashio::config 'operator_url')"
CLIENT_ID="$(bashio::config 'client_id')"
REGISTRATION_TOKEN="$(bashio::config 'registration_token')"
MODE="$(bashio::config 'mode')"
ALLOW_MUTATIONS="$(bashio::config 'allow_mutations')"
LOCAL_DB_MODE="$(bashio::config 'local_db_mode')"
POLL_INTERVAL="$(bashio::config 'poll_interval')"

mkdir -p /data/.talkha_state /data/.talkhalokal_state

cat > /data/.talkha.env <<EOF
# Generated by TalkHa Client add-on
HA_URL=${HA_URL}
HA_TOKEN=${HA_TOKEN}
EOF

export TALKHA_OPERATOR_URL="${OPERATOR_URL}"
export TALKHA_CLIENT_ID="${CLIENT_ID}"
export TALKHA_REGISTRATION_TOKEN="${REGISTRATION_TOKEN}"
export TALKHA_MODE="${MODE}"
export TALKHA_ALLOW_MUTATIONS="${ALLOW_MUTATIONS}"
export TALKHA_LOCAL_DB_MODE="${LOCAL_DB_MODE}"
export TALKHA_ADDON_MODE="1"
export TALKHA_LOCAL_DB="1"
export TALKHA_POLL_INTERVAL="${POLL_INTERVAL}"

exec python3 /opt/talkha/app/agent.py
